# deploy/systemd/legislative-reviews.service
[Unit]
Description=Legislative review pipeline
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/LegislativeReviews
EnvironmentFile=/opt/LegislativeReviews/.env
ExecStart=/opt/LegislativeReviews/.venv/bin/python scripts/review_worker_daemon.py
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
